# server/vino-server.service.in
[Unit]
Description=Vino VNC server

[Service]
Type=dbus
BusName=org.gnome.Vino
ExecStart=@libexecdir@/vino-server
Restart=on-failure
